# docker/supervisord.conf
[program:php-fpm]
command=/usr/local/sbin/php-fpm -F -O
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; Heratio doesn't ship a queue worker today; add one here when needed:
;   [program:queue-worker]
;   command=php /var/www/heratio/artisan queue:work --tries=3 --sleep=3
