# tools/supervisord/heratio-queue-worker.conf
; -----------------------------------------------------------------------------
; Heratio queue worker — supervisord program block
; -----------------------------------------------------------------------------
;
; Install:
;   sudo apt install supervisor
;   sudo cp /usr/share/nginx/heratio/tools/supervisord/heratio-queue-worker.conf \
;          /etc/supervisor/conf.d/heratio-queue-worker.conf
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start heratio-queue-worker:*
;
; Verify:
;   sudo supervisorctl status heratio-queue-worker:*
;   tail -f /var/log/heratio-queue-worker-00.log
;
; See: docs/queue-worker-deployment.md
; -----------------------------------------------------------------------------

[program:heratio-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php8.3 -d memory_limit=512M /usr/share/nginx/heratio/artisan queue:work --queue=default,integrations --sleep=3 --tries=3 --timeout=60 --memory=512
directory=/usr/share/nginx/heratio
user=www-data
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/heratio-queue-worker-%(process_num)02d.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=60
stopsignal=TERM
startsecs=5
startretries=3
priority=30
